# setup/etc/systemd/system/antizapret-update.service
[Unit]
Description=AntiZapret VPN update
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/root/antizapret
ExecStart=/bin/bash -c 'cat /root/antizapret/doall.sh | bash'
Restart=on-failure
RestartSec=1h
TimeoutSec=1h
SuccessExitStatus=0

[Install]
WantedBy=multi-user.target