# setup/etc/systemd/system/antizapret-update.timer
[Unit]
Description=antizapret-vpn update every 24h

[Timer]
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target