# setup/etc/systemd/system/openvpn-server@.service.d/override.conf
[Service]
LimitNPROC=infinity
LimitNOFILE=65535
StandardOutput=null
StandardError=null