# setup/etc/systemd/system/openvpn-server@.service.d/override.conf
[Unit]
Before=antizapret.service

[Service]
LimitNPROC=infinity
TasksMax=infinity
LimitNOFILE=infinity