# setup/etc/systemd/system/kresd@.service.d/override.conf
[Service]
LimitNPROC=infinity
LimitNOFILE=524288