# setup/etc/systemd/system/kresd@.service.d/override.conf
[Service]
LimitNPROC=infinity

[Unit]
Wants=network-online.target
After=network-online.target
Before=nss-lookup.target
Wants=nss-lookup.target