# setup/etc/systemd/system/wg-quick@.service.d/override.conf
[Service]
LimitNPROC=infinity
LimitNOFILE=65535

[Unit]
ConditionPathExists=/etc/wireguard/vpn.conf
ConditionPathExists=/etc/wireguard/antizapret.conf